# baler/etc/systemd/balerd.master.service.in
# README
# ======
# Please copy this file to /etc/systemd/system/ before
# applying any of the following commands.
#
# To make ldmsd starts automatically at system startup time:
#   systemctl enable balerd.master.service
#
# To disable automatic startup:
#   systemctl disable balerd.master.service
#
# To manually start the service:
#   systemctl start balerd.master.service
#
# To check status:
#   systemctl status balerd.master.service
#
# To stop the service:
#   sytemctl stop balerd.master.service

[Unit]
Description = Baler Daemon - Master Mode
Documentation = http://ovis.ca.sandia.gov

[Service]
Type = forking
EnvironmentFile = @sysconfdir@/baler/balerd.master.env
Environment = "HOSTNAME=$(hostname)"
ExecStart = @sbindir@/balerd \
		-s ${BALERD_STORE} \
		-C ${BALERD_PLUGIN_CONFIG} \
		-x sock \
		-p ${BALERD_MASTER_PORT} \
		-v ${BALERD_LOG_LEVEL} \
		-I ${BALERD_IN_Q_WORKERS} \
		-O ${BALERD_OUT_Q_WORKERS} \
		-m master

[Install]
WantedBy = default.target
