# baler/etc/systemd/bhttpd.service.in
# README
# ======
# Please copy this file to /etc/systemd/system/ before
# applying any of the following commands.
#
# To make bhttpd starts automatically at system startup time:
#   systemctl enable bhttpd.service
#
# To disable automatic startup:
#   systemctl disable bhttpd.service
#
# To manually start the service:
#   systemctl start bhttpd.service
#
# To check status:
#   systemctl status bhttpd.service
#
# To stop the service:
#   sytemctl stop bhttpd.service

[Unit]
Description = Baler HTTP Daemon
Documentation = http://ovis.ca.sandia.gov

[Service]
Type = forking
EnvironmentFile = @sysconfdir@/baler/bhttpd.env
Environment = "HOSTNAME=$(hostname)"
ExecStart = @sbindir@/bhttpd \
		-s ${BHTTPD_STORE} \
		-p ${BHTTPD_PORT} \
		-w ${BHTTPD_WORKERS} \
		-v ${BHTTPD_LOG_LEVEL}


[Install]
WantedBy = default.target
