# ldms/etc/systemd/ldmsd.aggregator.service.in
[Unit]
Description = LDMS Daemon
Documentation = http://ovis.ca.sandia.gov

[Service]
Type = forking
LimitNOFILE = ${LDMSD_ULIMIT_NOFILE}
EnvironmentFile = @sysconfdir@/ldms/ldmsd.aggregator.env
Environment = HOSTNAME=%H
ExecStartPre = /bin/mkdir -p @localstatedir@/run/ldmsd
ExecStart = @sbindir@/ldmsd \
		-x ${LDMSD_XPRT}:${LDMSD_PORT} \
		-c ${LDMSD_PLUGIN_CONFIG_FILE} \
		-a ${LDMSD_AUTH_PLUGIN} \
		-v ${LDMSD_VERBOSE} \
		-m ${LDMSD_MEM} \
		$LDMSD_LOG_OPTION \
		-P ${LDMSD_NUM_THREADS} \
		-r @localstatedir@/run/ldmsd/aggregator.pid

[Install]
WantedBy = default.target
