# ldms/etc/systemd/ldmsd.aggregator.service.in
# README
# ======
# Please copy this file and ldmsd.aggregator.env to /etc/systemd/system/ before
# applying any of the following commands.
#
# To make ldmsd starts automatically at system startup time:
#   systemctl enable ldmsd
#
# To disable automatic startup:
#   systemctl disable ldmsd
#
# To manually start ldmsd service:
#   systemctl start ldmsd
#
# To check status:
#   systemctl status ldmsd
#
# To stop the service:
#   sytemctl stop ldmsd
#

[Unit]
Description = LDMS Daemon
Documentation = http://ovis.ca.sandia.gov

[Service]
Type = forking
EnvironmentFile = @sysconfdir@/ldms/ldmsd.aggregator.env
Environment = HOSTNAME=%H
ExecStartPre = /bin/mkdir -p @localstatedir@/run/ldmsd
ExecStart = @sbindir@/ldmsd \
		-x ${LDMSD_XPRT}:${LDMSD_PORT} \
		-c ${LDMSD_PLUGIN_CONFIG_FILE} \
		-a ${LDMSD_AUTH_PLUGIN} \
		$LDMSD_AUTH_OPTS \
		-r @localstatedir@/run/ldmsd/aggregator.pid

[Install]
WantedBy = default.target
